# scripts/com.langbai.portfolio-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.langbai.portfolio-autopush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/langbai/virtual-venus/scripts/auto-push.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/langbai/virtual-venus</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/langbai/Library/Logs/portfolio-autopush/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/langbai/Library/Logs/portfolio-autopush/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/langbai</string>
    </dict>
</dict>
</plist>
